# Install/files/nmsprime-websockets.ini
[program:websockets]
command=/opt/rh/rh-php73/root/usr/bin/php /var/www/nmsprime/artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
user=apache
